# roles/cruft/templates/systemd.timer.j2
[Unit]
Description=Run Monad cruft service every hour

[Timer]
OnCalendar=hourly
Unit={{ service_identifier }}.service
Persistent=true

[Install]
WantedBy=multi-user.target
